# .config/autostart/1password.desktop
[Desktop Entry]
Name=1Password
Exec=/opt/1Password/1password %U --ozone-platform-hint=auto --silent
Terminal=false
Type=Application
Icon=1password
