# portables/forward-ports/mkosi.extra/usr/lib/systemd/system/forward-ports-http.service
[Unit]
Requires=forward-ports-http.socket
After=forward-ports-http.socket

[Service]
Type=notify
ExecStart=/usr/lib/systemd/systemd-socket-proxyd 127.0.0.1:1080
PrivateTmp=yes
